# supervisor/rootdb-reverb_server.conf
[program:rootdb-websocket_server]
process_name=%(program_name)s_%(process_num)02d
command=php /path/to/api/artisan reverb:start
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=<www-data or httpd>
numprocs=1
redirect_stderr=true
stdout_logfile=/path/to/api/storage/logs/websocket.log
stopwaitsecs=3600
