# supervisor/rootdb-cron_scheduler.conf
[program:rootdb-cron_scheduler]
process_name=%(program_name)s_%(process_num)02d
command=php /path/to/api/artisan schedule:run -q && exec /usr/bin/sleep 60
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=<www-data or httpd>
numprocs=1
redirect_stderr=true
