# docker-dev/supervisor/supervisor.d/rootdb-websocket_server.ini
[program:rootdb-websocket_server]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/api/artisan websockets:serve
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=rootdb
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/api/storage/logs/websocket.log
stopwaitsecs=3600
