# docker-dev/supervisor/supervisor.d/rootdb-cron-scheduler.ini
[program:rootdb-cron_scheduler]
process_name=%(program_name)s_%(process_num)02d
command=sh -c 'php /var/www/api/artisan schedule:run -q && exec /usr/bin/sleep 60'
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=rootdb
numprocs=1
redirect_stderr=true
