# docker-dev/supervisor/supervisor.d/rootdb-reverb_server.ini
[program:rootdb-reverb_server]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/api/artisan reverb:start
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=rootdb
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/api/storage/logs/reverb_websocket.log
stopwaitsecs=3600
